# roles/ddns/files/ddns.timer
[Unit]
Description=DDNS timer

[Timer]
OnStartupSec=30
OnUnitActiveSec=5min

[Install]
WantedBy=default.target
