# roles/mail/getmail/files/getmail.timer
[Unit]
Description=getmail timer

[Timer]
OnStartupSec=30
OnUnitActiveSec=3min
Unit=getmail.service